# hid-bridge/com.aislopdesk.hid-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchDaemon for aislopdesk-hid-bridge (runs as ROOT, since the Karabiner virtual-HID report socket is
  root-only). Lets videohostd type into SecurityAgent password dialogs over localhost UDP :9100.

  Install (one-time, needs sudo):
    bash hid-bridge/build.sh
    sudo cp hid-bridge/build/aislopdesk-hid-bridge /usr/local/bin/aislopdesk-hid-bridge
    sudo cp hid-bridge/com.aislopdesk.hid-bridge.plist /Library/LaunchDaemons/
    sudo launchctl load -w /Library/LaunchDaemons/com.aislopdesk.hid-bridge.plist
  Then launch videohostd with AISLOPDESK_VIRTUAL_HID=1.
  Uninstall: sudo launchctl unload -w /Library/LaunchDaemons/com.aislopdesk.hid-bridge.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aislopdesk.hid-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/aislopdesk-hid-bridge</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/aislopdesk-hid-bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/aislopdesk-hid-bridge.log</string>
</dict>
</plist>
